# s6-rc.d/postgres/run
#!/command/with-contenv sh
# Create subdirectories that may not exist in a fresh volume mount
mkdir -p /data/postgres /data/jobs /data/wd14-models /run/postgresql

# Initialize data directory if empty
if [ ! -f /data/postgres/PG_VERSION ]; then
    chown postgres:postgres /data/postgres
    su -s /bin/sh postgres -c "initdb -D /data/postgres"
fi

# Ensure ownership (handles pre-existing volumes with wrong perms)
chown -R postgres:postgres /data/postgres /run/postgresql

exec s6-setuidgid postgres postgres \
    -D /data/postgres \
    -c listen_addresses=127.0.0.1 \
    -c unix_socket_directories=/run/postgresql \
    -c logging_collector=off
